# deploy/k3s-lxc/k3s-lxc-fix.service
[Unit]
Description=K3s LXC Fixes
Before=k3s.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash -c "mount --make-rshared / && [ -e /dev/kmsg ] || ln -s /dev/console /dev/kmsg"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
